# scripts/data/rss-cache/_news/Microsoft_Security_Blog/httpswwwmicrosoftcomen-ussecurityblog20260417domain-compromise-predictive-shielding-shut-down-latera.json
{
  "Description": "Domain compromise accelerates fast. Predictive shielding slowed it down. This real-world attack shows how exposure-based containment stopped credential abuse and broke the threat actor's momentum.\n\nThe post [Containing a domain compromise: How predictive shielding shut down lateral movement](https://www.microsoft.com/en-us/security/blog/2026/04/17/domain-compromise-predictive-shielding-shut-down-lateral-movement/) appeared first on [Microsoft Security Blog](https://www.microsoft.com/en-us/security/blog).",
  "ProcessedDate": "2026-04-17 16:25:20",
  "Author": "Microsoft Defender Security Research Team",
  "OutputDir": "_news",
  "PubDate": "2026-04-17T14:51:01+00:00",
  "Link": "https://www.microsoft.com/en-us/security/blog/2026/04/17/domain-compromise-predictive-shielding-shut-down-lateral-movement/",
  "FeedUrl": "https://www.microsoft.com/en-us/security/blog/feed/",
  "Tags": [
    "Credential theft",
    "Domain compromise",
    "Lateral movement"
  ],
  "FeedName": "Microsoft Security Blog",
  "EnhancedContent": "In identity-based attack campaigns, any initial access activity can turn an already serious intrusion into a critical incident once it allows a threat actor to obtain domain-administration rights. At that point, the attacker effectively controls the Active Directory domain: they can change group memberships and Access Control Lists (ACLs), mint Kerberos tickets, replicate directory secrets, and push policy through mechanisms like Group Policy Objects (GPOs), among others.\n\nWhat makes domain compromise especially challenging is how quickly it could happen: in many real-world cases, domain-level credentials are compromised immediately following the very first access, and once these credentials are exposed, they’re often abused immediately, well before defenders can fully scope what happened. Apart from this speed gap, responding to this type of compromise could also prove difficult. For one, incident responders can’t just simply “turn off” domain controllers, service accounts, or identity infrastructure and core services without risking business continuity. In addition, because compromised credential artifacts can spread fast and be replayed to expand access, restoring the identity infrastructure back to a trusted state usually means taking steps (for example, [krbtgt rotation](https://learn.microsoft.com/en-us/windows-server/identity/ad-ds/manage/forest-recovery-guide/ad-forest-recovery-reset-the-krbtgt-password), [GPO cleanup](https://learn.microsoft.com/en-us/windows-server/identity/ad-ds/manage/group-policy/group-policy-backup-restore), and [ACL validation](https://learn.microsoft.com/en-us/windows/win32/ad/how-access-control-works-in-active-directory-domain-services)) that could take additional time and effort in an already high-pressure situation.\n\nThese challenges highlight the need for a more proactive approach in disrupting and containing credential-based attacks as they happen. Microsoft Defender’s predictive shielding capability in automatic attack disruption helps address this need. Its ability to predict where attacks will pivot next and apply just in time hardening actions to block credential abuse—including those targeting high-privilege accounts like domain admins—and lateral movement at near-real-time speed, shifting the advantageto the defenders.\n\nPreviously, we [discussed](https://www.microsoft.com/en-us/security/blog/2026/03/23/case-study-predictive-shielding-defender-stopped-gpo-based-ransomware-before-started/) how predictive shielding was able to disrupt a human-operated ransomware incident. In this blog post, we take a look at a real-world Active Directory domain compromise that illustrates the critical inflection point when a threat actor achieves domain -level control. We walk through the technical details of the incident to highlight attacker tradecraft, the operational challenges defenders face after domain compromise, and the value of proactive, exposure-based containment that predictive shielding provides.\n\n## Predictive shielding overview\n\n[Predictive shielding](https://learn.microsoft.com/en-us/defender-xdr/shield-predict-threats) is a capability in Microsoft Defender’s [automatic attack disruption](https://learn.microsoft.com/en-us/defender-xdr/automatic-attack-disruption) that helps stop the spread of identity-based attacks, before an attacker fully operationalizes stolen credentials. Instead of waiting for an account to be observed doing something malicious, predictive shielding focuses on moments when credentials are likely exposed: when Defender sees high-confidence signals of credential theft activity on a device, it can proactively restrict the accounts that might have been exposed there.\n\nEssentially, predictive shielding works as follows:\n\n- Defender detects post-breach activity strongly associated with credential exposure on a device.\n- It evaluates which high-privilege identities were likely exposed in that context.\n- It applies **containment** to those identities to reduce the attacker’s ability to pivot, limiting lateral movement paths and high-impact identity operations while the incident is being investigated and remediated. The intent is to close the “speed gap” where attackers can reuse newly exposed credentials faster than responders can scope, reset, and clean up.\n\nThis capability is available as an out-of-the-box enhancement for Microsoft Defender for Endpoint P2 customers who meet the Microsoft Defender prerequisites.\n\nThe following section revisits a real-world domain compromise that showcases how attack disruption and predictive shielding changed the outcome by acting on exposure, rather than just observed abuse. Interestingly, this case happened just as we’re rolling out the [predictive shielding](https://www.microsoft.com/en-us/security/blog/2026/03/23/case-study-predictive-shielding-defender-stopped-gpo-based-ransomware-before-started/?msockid=036f567c77ca6191158a408b767d60a9), so you can see the changes in both attacker tradecraft and the detection and response actions before and after this capability was deployed.\n\n## Attack chain overview\n\nIn June 2025, a public sector organization was targeted by a threat actor. This threat actor progressed methodically: initial exploitation, local escalation, directory reconnaissance, credential access, and expansion into Microsoft Exchange and identity infrastructure.\n\n![](https://www.microsoft.com/en-us/security/blog/wp-content/uploads/2026/04/image-30.webp)*Figure 1. Attack diagram of the domain compromise*.\n\n### Initial entry: Pre-domain compromise\n\nThe campaign began at the edge: a file-upload flaw in an internet-facing Internet Information Services (IIS) server was abused to plant and launch a web shell. The attacker then simultaneously performed various reconnaissance activities using the compromised account through the web shell and escalated their privileges to NT AUTHORITY\\SYSTEM by abusing a Potato-class token impersonation primitive (for example, BadPotato).\n\nThe discovery commands observed in the attack include the following example:\n\n![](https://www.microsoft.com/en-us/security/blog/wp-content/uploads/2026/04/image-31.webp)\n\nUsing the compromised IIS service account, the attacker attempted to reset the passwords of high-impact identities, a common technique used to gain control over accounts without performing credential dumping. The attacker also deployed Mimikatz to dump logon secrets (for example, MSV, LSASS, and SAM), harvesting credentials that are exposed on the device.\n\n![](https://www.microsoft.com/en-us/security/blog/wp-content/uploads/2026/04/image-32.webp) ![](https://www.microsoft.com/en-us/security/blog/wp-content/uploads/2026/04/image-33.webp)\n\nHad predictive shielding been released at this point, automated restrictions on exposed accounts could have stopped the intrusion before it expanded beyond the single-host foothold. However, at the time of the incident, this capability hasn’t been deployed to customers yet.\n\n**Key takeaway:** At this stage of an attack, it’s important to keep the containment host‑scoped. Defenders should prioritize blocking credential theft and stopping escalation before it reaches the identity infrastructure.\n\n### First pivot: Directory credential materialization and Exchange delegation\n\nWithin 24 hours, the attacker abused privileged accounts and remotely created a scheduled task on a domain controller. The task initiated NTDS snapshot activity and packaged the output using makecab.exe, enabling offline access to directory credential material that’s suitable for abusing credentials at scale:\n\n![](https://www.microsoft.com/en-us/security/blog/wp-content/uploads/2026/04/image-34.webp)\n\nBecause the first malicious action by the abused account already surfaced the entire Active Directory credentials, stopping its path for total domain compromise was no longer feasible.\n\nThe threat actor then planted a Godzilla web shell on Exchange Server, used a privileged context to enumerate accounts with *ApplicationImpersonation* role assignments, and granted full access to a delegated principal across mailboxes using *Add‑MailboxPermission*. This access allowed the threat actor to read and manipulate all mailbox contents.\n\nThe attack also used Impacket’s *atexec.py* to enumerate the role assignments remotely. Its use triggered the attack disruption capability in Defender, revoking the account sessions of an admin account and blocking it from further use.\n\n![](https://www.microsoft.com/en-us/security/blog/wp-content/uploads/2026/04/image-35.webp)\n\nFollowing the abused account’s disruption, the attacker attempted several additional actions, such as resetting the disrupted account’s and other accounts’ passwords. They also attempted to dump credentials of a Veeam backup device.\n\n**Key takeaway:** This pivot is a turning point. Once directory credentials and privileged delegation are in play, the scope and impact of an incident expand fast. Defenders should prioritize protecting domain controllers, privileged identities, and authentication paths.\n\n### Scale and speed: Tool return, spraying, and lateral movement\n\nWeeks later, the threat actor returned with an Impacket tooling (for example, secretsdump and PsExec) that resulted in repeated disruptions by Defender against the abused accounts that they used. These disruptions forced the attacker to pivot to other compromised accounts and exhaust their resources.\n\nFollowing Defender’s disruptions, the threat actor then launched a broad password spray from the initially compromised IIS server, unlocking access to at least 14 servers through password reuse. They also attempted remote credential dumping against a couple of domain controllers and an additional IIS server using multiple domain and service principals.\n\n**Key takeaway:** Even though automatic attack disruption acted right away, the attacker already possessed multiple credentials due to the previous large-scale credential dumping. This scenario showcases the race to detect and disrupt credential abuse and is the reason we’re introducing predictive shielding to preemptively disrupt exposed accounts at risk.\n\n### Predictive shielding breaks the chain: Exposure-centric containment\n\nIn the second phase of the attack, we activated predictive shielding. When exposure signals surfaced (for example, credential dumping attempts and replay from compromised hosts), automated containment blocked new sign-in attempts and interactive pivots not only for the abused accounts, but also for context-linked identities that are active on the same compromised surfaces.\n\nAttack disruption contained high-privileged principals to prevent these accounts from being abused. Crucially, when a high-tier Enterprise or Schema Admin credential was exposed, predictive shielding contained it pre-abuse, preventing what would normally become a catastrophic escalation.\n\n### Second pivot: Alternative paths to new credentials\n\nWith high-value identities pre-contained, the threat actor pivoted to exploiting Apache Tomcat servers. They compromised three Tomcat servers, dropped the Godzilla web shell, and launched the PowerShell-based *Invoke-Mimikatz* command to harvest additional credentials. At one point, the attacker operated under Schema Admin:\n\n![](https://www.microsoft.com/en-us/security/blog/wp-content/uploads/2026/04/image-36.webp)\n\nThey then used Impacket *WmiExec* to access Microsoft Entra Connect servers and attempt to extract Entra Connect synchronization credentials. The account used for this pivot was later contained, limiting further lateral movement.\n\n### Last attempts and shutdown\n\nIn the final phase of the attack, the threat actor attempted a full LSASS dump on a file sharing server using *comsvcs.dll MiniDump* under a domain user account, followed by additional NTDS activity:\n\n![](https://www.microsoft.com/en-us/security/blog/wp-content/uploads/2026/04/image-37.webp)\n\nAttack disruption in Defender repeatedly severed sessions and blocked new sign-ins made by the threat actor. On July 28, 2025, the attack campaign lost momentum and stopped.\n\n## How predictive shielding changed the outcome\n\nBefore compromising a domain, attackers are mostly constrained by the hosts they control. However, even a small set of exposed credentials could remove their constraints and give them broad access through privileged authentication and delegated pathways. The blast radius spreads fast, time pressure spikes, and containment decisions become riskier because identity infrastructure and high-privilege accounts are production dependencies.\n\nThe incident we revisited earlier almost followed a similar pattern. It unfolded while predictive shielding was still being launched, so the automated predictive containment capability only became active at the midway of the attack campaign. During the attack’s first stages, the threat actor had room to scale—they returned with new tooling, launched a broad password spray attack, and expanded access across multiple servers. They also attempted remote credential dumping against domain controllers and servers.\n\nWhen predictive shielding went live, it helped shift the story and we then saw the change of pace—instead of reacting to each newly abused account, the capability allowed Defender to act preemptively and turn credential theft attempts into blocked pivots. Defender was able to block new sign-ins and interactive pivots, not just for the single abused account, but also for context-linked identities that were active on the same compromised surfaces.\n\nWith high-value identities pre-contained, the adversary shifted tradecraft and chased other credential sources, but each of their subsequent attempts triggered targeted containment that limited their lateral reach until they lost momentum and stopped. How this incident concluded is the operational “tell” that containment is working, in that once privileged pivots get blocked, threat actors often hunt for alternate credential sources, and defenses must continue following the moving blast radius.\n\nAs predictive shielding matures, it will continue to expand its prediction logic and context-linked identities.\n\n## MITRE ATT&CK® techniques observed\n\nThe following table maps observed behaviors to ATT&CK®.\n\n*Tactics shown are per technique definition.*\n\n| **Tactic(s)** | **Technique ID** | **Technique name** | **Observed details** | | --- | --- | --- | --- | | Initial Access | T1190 | Exploit Public-Facing Application | Exploited a file-upload vulnerability in an IIS server to drop a web shell. | | Persistence | T1505.003 | Server Software Component: Web Shell | Deployed web shells for persistent access. | | Execution | T1059.001 | Command and Scripting Interpreter: PowerShell | Used PowerShell for Exchange role queries, mailbox permission changes, and Invoke-Mimikatz. | | Privilege Escalation | T1068 | Exploitation for Privilege Escalation | Used BadPotato to escalate to SYSTEM on an IIS server. | | Credential Access | T1003.001 | OS Credential Dumping: LSASS Memory | Dumped LSASS using Mimikatz and *comsvcs.dll MiniDump*. | | Credential Access | T1003.003 | OS Credential Dumping: NTDS | Performed NTDS-related activity using ntdsutil snapshot/IFM workflows on a domain controller. | | Execution; Persistence; Privilege Escalation | T1053.005 | Scheduled Task/Job: Scheduled Task | Created remote scheduled tasks to execute under SYSTEM on a domain controller. | | Discovery | T1087.002 | Account Discovery: Domain Account | Enumerated domain groups and accounts using net group and AD Explorer. | | Lateral Movement | T1021.002 | Remote Services: SMB/Windows Admin Shares | Used admin shares/SMB-backed tooling (for example, PsExec) for lateral movement. | | Lateral Movement | T1021.003 | Remote Services: Windows Remote Management | Used WmiExec against Microsoft Entra Connect servers. | | Credential Access | T1110.003 | Brute Force: Password Spraying | Performed password spraying leading to access across at least 14 servers. | | Collection | T1114.002 | Email Collection: Remote Email Collection | Expanded mailbox access broadly through impersonation or permission changes. | | Command and Control | T1071.001 | Application Layer Protocol: Web Protocols | Web shells communicated over HTTP/S. | | Defense Evasion | T1070.004 | Indicator Removal on Host: File Deletion | Used cleanup scripts (for example, *del.bat*) to remove dump artifacts. | | Persistence; Privilege Escalation | T1098 | Account Manipulation | Manipulated permissions and roles to expand access and sustain control. | | Credential Access | T1078 | Valid Accounts | Reused compromised service and domain accounts for access and lateral movement. |\n\n## Learn more\n\n**For more information about automatic attack disruption and predictive shielding, see the following Microsoft Learn articles:**\n\n- Check out our latest Ninja [show](https://youtu.be/u9nXRABIw1k?si=A6X9q1rHFqKPVtQb) showcasing how predictive shielding expands to identity centric attacks\n- [Automatic attack disruption in Microsoft Defender XDR](https://learn.microsoft.com/en-us/defender-xdr/automatic-attack-disruption)\n- [Predictive shielding in Microsoft Defender (Preview)](https://learn.microsoft.com/en-us/defender-xdr/shield-predict-threats)",
  "FeedLevelAuthor": "Microsoft Security Blog",
  "Title": "Containing a domain compromise: How predictive shielding shut down lateral movement"
}
